# files/usr/lib/systemd/user/ostree-update-notify.service
[Unit]
Description=Notify about staged OSTree updates

[Service]
Type=oneshot
ExecStart=/usr/libexec/ostree-update-notify.sh